# init/znapzend.conf
# Ubuntu upstart file at /etc/init/znapzend.conf

respawn

start on runlevel [2345]
stop on runlevel [06]

expect fork

exec /usr/local/bin/znapzend --daemonize
