# init/znapzend.upstart
# znapzend: upstart file at /etc/init/znapzend.conf                                                                                                                                                                                   

description "znapzend"

respawn

start on [2345]
stop on runlevel [06]

expect fork
reload signal HUP

exec /usr/local/bin/znapzend --daemonize
