# init/znapzend.service.in
[Unit]
Description=ZnapZend - ZFS Backup System
Documentation=man:znapzend
After=zfs-import-cache.service
After=zfs-import-scan.service

[Service]
EnvironmentFile=-/etc/default/znapzend
ExecStart=@BINDIR@/znapzend $ZNAPZENDOPTIONS
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
# might be neccessary on low power systems
# Nice=19
# IOSchedulingClass=2
# IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
